# .config/systemd/user/autobeet.service
[Unit]
Description=Beet Daemon

[Service]
ExecStart=/home/rodrigo/scripts/autobeet.sh
Restart=always

[Install]
WantedBy=default.target
