# systemd/user/mopidy.service
[Unit]
Description=Mopidy Music Player Daemon

[Service]
ExecStart=/home/rodrigo/.local/bin/mopidy

[Install]
WantedBy=default.target
